# scripts/launchd/com.openclaw.audit-cron.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.audit-cron</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/n8n/scripts/audit-cron.sh</string>
        <string>--audit-only</string>
    </array>

    <!-- Weekly Sunday at 03:00. Change --audit-only to --auto-fix for auto-remediation. -->
    <!-- If the Mac is asleep at 03:00, launchd runs at next wake. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Run as root for full system access -->
    <key>UserName</key>
    <string>root</string>

    <!-- Redirect stdout/stderr for debugging -->
    <key>StandardOutPath</key>
    <string>/opt/n8n/logs/audit/launchd-cron-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/n8n/logs/audit/launchd-cron-stderr.log</string>

    <!-- Do not restart on failure -->
    <key>KeepAlive</key>
    <false/>

    <!-- Background priority -->
    <key>Nice</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
